# systemd/constitutional-ai-backend.service
[Unit]
Description=Constitutional AI Backend
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
WorkingDirectory=/home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/backend
Environment="PATH=/home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ai-server/AN-FOR-NO-ASSHOLES/09_CONSTITUTIONAL-AI/backend/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=default.target

